# docs/systemd/flask-url-shortener.service
[Unit]
Description=Flask URL Shortener Application
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/flask-url-shortener
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
